# templates/cdk.master.auth-webhook.service
[Unit]
Description=CDK master auth webhook
After=network.target

[Service]
User=root
WorkingDirectory=/root/cdk
ExecStart={{ charm_dir }}/../.venv/bin/gunicorn \
    --bind {{ host }}:{{ port }} \
    --capture-output \
    --certfile ./server.crt \
    --disable-redirect-access-to-syslog \
    --error-logfile /root/cdk/audit/auth-webhook.log \
    --keyfile ./server.key \
    --log-level debug \
    --workers 2 \
    auth-webhook:app
Restart=always

[Install]
WantedBy=multi-user.target
